# obsidian/supervisord.conf
[supervisord]
nodaemon=true
user=obsidian
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_RESOLUTION)s
autorestart=true
priority=10

[program:fluxbox]
command=fluxbox
environment=DISPLAY=":1"
autorestart=true
priority=20

[program:x11vnc]
command=x11vnc -display :1 -forever -shared -nopw -rfbport %(ENV_VNC_PORT)s
autorestart=true
priority=30

[program:novnc]
command=websockify --web /usr/share/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autorestart=true
priority=40

[program:obsidian]
command=obsidian --no-sandbox --disable-gpu /vault
environment=DISPLAY=":1"
autorestart=true
priority=50
startsecs=3
